# initrd/busybox-1.32.1/_install-bpi/etc/inittab
::sysinit:/etc/init.d/rcS
ttyS0::respawn:/usr/sbin/chroot /mnt /sbin/getty "-L" "ttyS0" "115200" "xterm-256color"
