# initrd/busybox-1.32.1/_install-emu/etc/inittab
::sysinit:/etc/init.d/rcS
ttyAMA0::respawn:/usr/sbin/chroot /mnt /sbin/getty "-L" "ttyAMA0" "115200" "linux"
