# initrd/busybox-1.33-libra_spl/etc/inittab
::sysinit:/etc/init.d/rcS
ttymxc0::respawn:/sbin/getty "-L" "ttymxc0" "115200" "xterm-256color"
